# docker/cron/crontab
* * * * * cd /var/www/html && /usr/local/bin/php artisan zoo:make-turn >> /var/log/cron.log 2>&1
# Пустая строка в конце файла обязательна
